# chrome.ini
[program:chrome]
command = /usr/bin/chromedriver
autostart=true
autorestart=true
priority=50
startretries=3
stopwaitsecs=10
